# mi_app.service
[Unit] 
Description=Mi Script de Python de Ejemplo 
# Esta línea asegura que el script se ejecute después de que la red esté lista. 
# ¡Es crucial para scripts que necesitan conexión a internet o a la red local!
After=network.target

[Service]
# Usuario que ejecutará el script (es mejor no usar root si no es necesario) 
User=pi 
# Grupo del usuario 
Group=pi 
# Directorio de trabajo: aquí es donde se ejecutará el script. 
# Ayuda a que las rutas relativas en tu script funcionen correctamente.
WorkingDirectory=/home/pi/scripts 
# El comando para iniciar el script. 
#Usamos la ruta completa. 
ExecStart=/home/pi/scripts/mi_script.py 
# Reiniciar el servicio automáticamente si falla. 
Restart=always 

[Install] 
WantedBy=multi-user.target
